# deployment/systemd/stock-trading-orchestrator.service
[Unit]
Description=Stock Trading Orchestrator Service (Master Controller)
Documentation=https://github.com/sangrokhan/ko-stock-filter
After=network.target stock-trading-postgres.service stock-trading-redis.service
After=stock-trading-data-collector.service stock-trading-indicator-calculator.service
After=stock-trading-stock-screener.service stock-trading-trading-engine.service
After=stock-trading-risk-manager.service
Requires=stock-trading-postgres.service stock-trading-redis.service
Wants=network-online.target

[Service]
Type=simple
User=stocktrading
Group=stocktrading
WorkingDirectory=/opt/stock-trading/services/orchestrator

# Environment
Environment="DATABASE_URL=postgresql://stock_user:stock_password@localhost:5432/stock_trading"
Environment="REDIS_HOST=localhost"
Environment="REDIS_PORT=6379"
Environment="LOG_LEVEL=INFO"
Environment="DATA_COLLECTOR_URL=http://localhost:8001"
Environment="INDICATOR_CALCULATOR_URL=http://localhost:8002"
Environment="STOCK_SCREENER_URL=http://localhost:8003"
Environment="TRADING_ENGINE_URL=http://localhost:8004"
Environment="RISK_MANAGER_URL=http://localhost:8005"
Environment="TZ=Asia/Seoul"
EnvironmentFile=-/etc/stock-trading/orchestrator.env

# Execution
ExecStart=/opt/stock-trading/venv/bin/python main.py
ExecReload=/bin/kill -HUP $MAINPID

# Restart policy
Restart=on-failure
RestartSec=10s

# Resource limits
LimitNOFILE=8192
LimitNPROC=2048

# Timeout settings
TimeoutStartSec=60
TimeoutStopSec=30

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stock-orchestrator

[Install]
WantedBy=multi-user.target
